# integralstor2/install/conf-files/services/first-boot.service
[Unit]
Description=Actions to be performed during the first boot
After=network.target

[Service]
Type=simple
ExecStart=/opt/integralstor/integralstor2/scripts/shell/first_boot.sh
RemainAfterExit=false

[Install]
WantedBy=multi-user.target
